# systemd/echo-restrict.socket
[Unit]
Description=test
[Socket]
ListenStream=127.0.0.1:9000

[Install]
WantedBy=default.target
